# stage-openmower/30-docker/files/etc/systemd/system/docker-preload-images.service
[Unit]
Description=Preload bundled Docker images (*.tar / *.tar.gz) into local Docker
After=docker.service
Requires=docker.service
Before=dockge.service webterminal.service
ConditionDirectoryNotEmpty=/opt/docker-images

[Service]
Type=oneshot
ExecStart=/bin/sh -e -c '\
dir=/opt/docker-images; \
# Decompress any .tar.gz \
for gz in "$dir"/*.tar.gz; do \
  [ -e "$gz" ] || break; \
  gunzip -f "$gz"; \
done; \
# Load all .tar and delete on success \
for tar in "$dir"/*.tar; do \
  [ -e "$tar" ] || break; \
  if docker load -i "$tar"; then rm -f "$tar"; fi; \
done \
'

[Install]
WantedBy=multi-user.target
